# common/.config/systemd/user/opencode-cursor-sync.service
[Unit]
Description=Refresh open-cursor (cursor-acp) model catalog for OpenCode
Documentation=https://github.com/Nomadcxx/opencode-cursor
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# Writes ~/.config/opencode/open-cursor.generated.json (gitignored), which OpenCode
# auto-loads via OPENCODE_CONFIG (set in ~/.config/environment.d/opencode.conf).
# npx is resolved through the mise shim so PATH does not need to be inherited.
ExecStart=%h/.local/share/mise/shims/npx -y @rama_nigg/open-cursor@latest sync-models --variants --compact --config %h/.config/opencode/open-cursor.generated.json --no-backup
LogLevelMax=warning
