# common/.config/systemd/user/proton-pass-cli-autologin.service
[Unit]
Description=Proton Pass CLI auto-login, SSH agent, and health monitor
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
KillMode=control-group
TimeoutStopSec=5
Environment=SSH_AGENT_SOCKET=%h/.ssh/proton-pass-agent.sock
Environment=SSH_AGENT_WAIT_SECONDS=120
Environment=MAX_WIFI_WAIT_SECONDS=45
Environment=MAX_NETWORK_WAIT_SECONDS=30
Environment=NETWORK_CHECK_HOST=1.1.1.1
Environment=LOGIN_VERIFY_WAIT_SECONDS=20
Environment=LOGIN_RETRY_INTERVAL=30
Environment=KEYRING_WAIT_SECONDS=90
Environment=HEALTHCHECK_INTERVAL=300
Environment=PROTON_LOGIN=%h/Applications/proton-login
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=%h/.local/bin/proton-pass-cli-autologin-monitor
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
