# htpc/.config/systemd/user/steam-controller-tv-switch.service
[Unit]
Description=Switch the HTPC to the TV when the Steam Controller connects
Documentation=man:udevadm(8)

[Service]
Type=simple
EnvironmentFile=%h/.config/htpc-tv-switch/env
ExecStart=%h/.local/bin/steam-controller-tv-switch-monitor
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=default.target
