# common/crontab
@reboot flatpak run --command=goldwarden com.quexten.Goldwarden daemonize
